# files/etc/crontabs/root
#min hour day month dayofweek command
*    *    *   *     *         /usr/local/bin/plan-b-check-tor

